# Alphonse/alphonse.service.in
[Unit]
Description=Alphonse phone handler
After=multi-user.target

[Service]
Type=simple

User=root
Group=root

WorkingDirectory={{ALPHONSE_DATA_DIR}}
ExecStart={{ALPHONSE_DEP_PYTHON_BIN}} {{ALPHONSE_BIN_DIR}}alphonse.py &> /dev/null

Restart=on-failure

[Install]
WantedBy=multi-user.target

